# packaging/client/rootfs/DEBIAN/prerm
#!/bin/sh
set -e
systemctl disable --now tunneler-client.service || true
exit 0