# packaging/server/rootfs/DEBIAN/prerm
#!/bin/sh
set -e
systemctl disable --now tunneler-server.service || true
exit 0
